# qmwc.service
[Unit]
Description=QuietMisdreavus Wallpaper Cycler

[Service]
Type=oneshot
ExecStart=%h/.cargo/bin/qmwc -vvv
